# crontab
#crontab -e

*/1 * * * * /home/bitnami/htdocs/auto-pull/git-pull-requests/fetch.sh 
